# supervisord-hs.conf
[program:hydroserver]
command=/final_run.sh
numprocs=1
autostart=true
autorestart=false
stdout_logfile = /mylogs/hs.log
redirect_stderr = true
startsecs = 0
