# supervisord-hs.conf
[program:hs]
command=/final_run.sh
numprocs=1
autostart=true
autorestart=true
